# etc/led_driver.ini
[program:worker]
command=/home/pi/.local/bin/celery -A pi_led worker --loglevel=warning
environment=DEBUG="false"
autorestart=true
redirect_stderr=true
stdout_logfile=/home/pi/pi_led/logs/worker.log
user=pi
directory=/home/pi/pi_led

[program:beat]
command=/home/pi/.local/bin/celery -A pi_led beat --loglevel=warning --scheduler=django_celery_beat.schedulers:DatabaseScheduler
environment=DEBUG="false"
autorestart=true
redirect_stderr=true
stdout_logfile=/home/pi/pi_led/logs/beat.log
user=pi
directory=/home/pi/pi_led

[program:django]
command=python3 manage.py runserver --noreload 0.0.0.0:8000
environment=DEBUG="false"
autorestart=true
redirect_stderr=true
stdout_logfile=/home/pi/pi_led/logs/supervisor.log
user=pi
directory=/home/pi/pi_led

[program:server]
command=python3 -m http.server --bind 0.0.0.0 8100
autorestart=true
redirect_stderr=true
stdout_logfile=/home/pi/pi_led/logs/server.log
user=pi
directory=/home/pi/pi_led/mobile/www
